# crontab
0 5 * * * /bin/bash /usr/local/bin/entrypoint.sh >> /var/log/cron.log 2>&1
