# crontab
0 5 * * * /bin/bash /opt/duplicity/duplicity_script.sh >> /var/log/cron.log 2>&1
